# packaging/systemd/barricade-beat.service
[Unit]
Description=Barricade Celery Beat Scheduler
After=network-online.target redis.service barricade-worker.service
Requires=redis.service
PartOf=barricade.target

[Service]
Type=simple
User=barricade
Group=barricade
WorkingDirectory=/usr/lib/barricade/backend
EnvironmentFile=/etc/barricade/barricade.env
Environment=PATH=/usr/lib/barricade/venv/bin:/usr/bin:/bin
ExecStart=/usr/lib/barricade/venv/bin/celery -A app.tasks beat --scheduler redbeat.RedBeatScheduler --loglevel=info
Restart=on-failure
RestartSec=5s
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/barricade /var/log/barricade /tmp /dev/shm

[Install]
WantedBy=barricade.target
