# packaging/systemd/barricade.service
[Unit]
Description=Barricade
After=network-online.target postgresql.service redis.service
Requires=postgresql.service redis.service

[Service]
Type=simple
User=barricade
Group=barricade
WorkingDirectory=/usr/lib/barricade/backend
Environment=PATH=/usr/lib/barricade/venv/bin:/usr/bin:/bin
Environment=BARRICADE_CONFIG=/etc/barricade/barricade.toml
ExecStart=/usr/lib/barricade/venv/bin/python -m app
Restart=on-failure
RestartSec=5s
TimeoutStopSec=65
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/barricade /var/log/barricade /tmp /dev/shm

[Install]
WantedBy=multi-user.target
