# packaging/systemd/labdog.service
[Unit]
Description=LabDog
Documentation=https://open-labdog.github.io/labdog/
# Order after network and (optional) local DB/Redis. We do NOT Require= them
# because PostgreSQL and Redis are commonly hosted off-box in production.
# The application retries DB/Redis connections at startup.
After=network-online.target postgresql.service redis.service
Wants=network-online.target

[Service]
Type=simple
User=labdog
Group=labdog
WorkingDirectory=/usr/lib/labdog/backend
Environment=PATH=/usr/lib/labdog/venv/bin:/usr/bin:/bin
Environment=LABDOG_CONFIG=/etc/labdog/labdog.toml
ExecStart=/usr/lib/labdog/venv/bin/python -m app
Restart=on-failure
RestartSec=5s
TimeoutStopSec=65

# Hardening
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
# /dev/shm is needed for the per-run SSH key tmpfs.
ReadWritePaths=/var/lib/labdog /var/log/labdog /run/labdog /dev/shm

[Install]
WantedBy=multi-user.target
